# src/DEBIAN/postinst
#!/bin/bash

if [ "$1" == "configure" ]; then
	update-rc.d jailcontrol defaults
fi
